# templates/99-usb-roms.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z][0-9]", \
    ENV{ID_FS_TYPE}!="", \
    RUN+="{{SCRIPTS_DIR}}/usb_roms.sh %k"
